# indird/indirdwake@.service
[Service]
StartLimitInterval=0
ExecStart=/usr/local/bin/indird %i wakeup
