# indird/indirdwake@.service
[Service]
ExecStart=/usr/local/bin/indird %i wakeup
