# indird/indirdwake@.path
[Unit] 
Description=InDirD %i path wakeup
Requires=indird@%i.service
After=indird@%i.service

[Path]
DirectoryNotEmpty=/run/indird/%i_path

[Install]
WantedBy=multi-user.target
