# indird/indirdwake@.service
[Unit]
Description=InDirD %i wakeup

[Service]
ExecStart=/usr/local/bin/indird %i wakeup
